# etc/service/log-hello-world/run
#!/bin/sh

exec /bin/sh -c "while true; do logger hello world; sleep 1; done;"
